# components/notebook-images/r-studio-experimental-no-apache/root/etc/services.d/rstudio/finish
#!/bin/bash
exec \
	s6-setuidgid jovyan \
                rstudio-server stop
