# components/notebook-images/code-server-experimental-no-apache/root/etc/services.d/code-server/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid jovyan \
		/usr/local/bin/code-server \
			--bind-addr 0.0.0.0:8888 \
			--disable-telemetry \
			--auth none
