# components/example-notebook-servers/rstudio/s6/services.d/rstudio/finish
#!/command/with-contenv bash
exec rstudio-server stop
