# components/notebook-images/r-studio-experimental/root/etc/services.d/apache2/run
#!/usr/bin/with-contenv bash

exec /usr/sbin/apachectl -DFOREGROUND
