# components/example-notebook-servers/codeserver/root/etc/services.d/code-server/run
#!/usr/bin/with-contenv bash

export SHELL='/bin/bash'
exec \
	/usr/local/bin/code-server \
		--bind-addr 0.0.0.0:8888 \
		--disable-telemetry \
		--auth none
