# components/example-notebook-servers/codeserver/s6/services.d/code-server/run
#!/usr/bin/with-contenv bash
exec /usr/bin/code-server \
  --bind-addr 0.0.0.0:8888 \
  --disable-telemetry \
  --auth none
